# php/cronjob/crontab
* * * * * root php /var/www/html/index.php >> /var/log/cron.log 2>&1
